# install_files/securedrop-config-focal/etc/systemd/system/apt-daily.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 00/3:00
RandomizedDelaySec=1h
Persistent=true
