# install_files/ansible-base/roles/build-securedrop-app-code-deb-pkg/files/cron.daily
#!/bin/bash

cd /var/www/securedrop || exit 1
./manage.py how-many-submissions-today
