# install_files/ansible-base/roles/build-ossec-deb-pkg/files/ossec-agent.service
[Unit]
Description=OSSEC service for agent

[Service]
Type=forking
ExecStart=/var/ossec/bin/ossec-control start
ExecStop=/var/ossec/bin/ossec-control stop
RemainAfterExit=True
SuccessExitStatus=1

[Install]
WantedBy=multi-user.target
