# securedrop/debian/app-code/lib/systemd/system/securedrop-clean-tmp.timer
[Unit]
Description=Run clean-tmp daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
