# securedrop/debian/app-code/lib/systemd/system/securedrop-remove-pending-sources.timer
[Unit]
Description=remove pending SecureDrop sources daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
